# debian/bloom.service
[Unit]
Description=Bloom HTTP REST API caching middleware
After=network.target

[Service]
Type=forking
User=bloom
Group=bloom
ExecStart=/usr/bin/bloom -c /etc/bloom.cfg
Restart=on-failure

[Install]
WantedBy=multi-user.target
